# etc/init.d/kano-settings
#!/bin/sh

### BEGIN INIT INFO
# Provides:         kano-settings
# Required-Start:
# Required-Stop:
# Default-Start:    2
# Default-Stop:
### END INIT INFO

/usr/bin/kano-settings-onboot
