# etc/init.d/kano-hdmi
#!/bin/sh

### BEGIN INIT INFO
# Provides:         kano-hdmi
# Required-Start:   mountkernfs
# Required-Stop:
# Default-Start:    2
# Default-Stop:
### END INIT INFO

/usr/bin/set-kano-hdmi
